# .config/systemd/user/gnome-openbox.target
[Unit]
Description=Kiwi Desktop Environment
Wants=gnome.target
Wants=openbox.target
AllowIsolate=true

[Install]
Alias=default.target
