# .config/systemd/user/gnome-session.service
[Unit]
Description=The GNOME session manager
Requires=xorg.target
After=xorg.target

[Service]
ExecStart=/bin/gnome-session
Restart=on-failure
