# .config/systemd/user/gnome.target
[Unit]
Description=The GNOME Desktop Environment
Wants=gnome-session.service
AllowIsolate=true
