# .config/systemd/user/kiwi-desktop.target
[Unit]
Description=Kiwi Desktop Environment
Wants=xorg.target
Requires=dbus.socket
AllowIsolate=true
