# .config/systemd/user/panel.target

[Unit]
Description=Desktop Panel
Wants=xorg.target
Requires=dbus.socket
